# roles/nginx/templates/nginx_app.path.j2
[Unit] 
Wants= nginx_reload.service 

[Path] 
PathChanged= {{ nginx_dir_path }}/conf/default.conf 
Unit=nginx_reload.service