# roles/nginx/templates/nginx_restart.service.j2
[Service]  
Type= simple  
ExecStart= docker-compose -f {{ nginx_dir_path }}/docker-compose.yml restart
