# roles/nginx/templates/nginx_conf.path.j2
[Unit] 
Wants= nginx_restart.service 

[Path] 
PathChanged= {{ nginx_dir_path }}/conf/nginx.conf 
